# .config/systemd/user/duckdns.timer
[Unit]
Description=update duckdns address

[Timer]
OnCalendar=hourly
Unit=duckdns.service

[Install]
WantedBy=timers.target
