# .config/systemd/user/duckdns.timer
[Unit]
Description=update duckdns address

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
