# .config/systemd/user/spotifyd.service.d/override.conf
[Service]
ExecStart=
ExecStart=/bin/sh -lc "/usr/bin/spotifyd --no-daemon"
