# .config/systemd/user/redshift.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/redshift -v