# .config/systemd/user/gammastep.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/gammastep -v
